# scripts/systemd/discord-bridge.service
[Unit]
Description=Discord bridge bot
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/vinted-ai-cloud
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/vinted-ai-cloud/.venv/bin/python tools/discord_bridge_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
